# rootfs/usr/lib/systemd/user/hypridle.service.d/10-arch-bootc.conf
[Service]
ExecStart=
ExecStart=/usr/bin/hypridle --config /usr/share/hypr/hypridle.conf
