# rootfs/usr/lib/systemd/user/hypridle.service.d/10-arch-bootc.conf
[Service]
Environment=XDG_CONFIG_DIRS=/usr/share:/usr/share/xdg:/etc/xdg
ExecStart=
ExecStart=/usr/bin/hypridle
