# rootfs/etc/systemd/system/python3-validity.service.d/device-presence.conf
[Unit]
# Skip the Validity/Synaptics fingerprint daemon on systems without a matching
# USB fingerprint reader. In VMs this service otherwise fails at boot and
# obscures real regressions.

[Service]
ExecCondition=/usr/bin/bash -c '\
shopt -s nullglob; \
for dev in /sys/bus/usb/devices/*; do \
  [ -r "$dev/idVendor" ] || continue; \
  vendor=$(cat "$dev/idVendor" 2>/dev/null || true); \
  product=$(cat "$dev/idProduct" 2>/dev/null || true); \
  case "$vendor:$product" in \
    138a:*|06cb:009a) exit 0 ;; \
  esac; \
done; \
exit 1'
