# rootfs/usr/lib/systemd/user/waybar.service.d/10-arch-bootc.conf
[Service]
ExecStart=
ExecStart=/usr/bin/waybar --config /usr/share/waybar/arch-bootc/config.jsonc --style /usr/share/waybar/arch-bootc/style.css
